# semi-agentic.service.template
# Systemd Service Template for Intel NUC
# 
# To use this:
# 1. Edit the paths and username below to match your system
# 2. Copy to /etc/systemd/system/semi-agentic.service
#    sudo cp semi-agentic.service.template /etc/systemd/system/semi-agentic.service
# 3. Edit the copied file with sudo
#    sudo nano /etc/systemd/system/semi-agentic.service
# 4. Enable and start:
#    sudo systemctl daemon-reload
#    sudo systemctl enable semi-agentic.service
#    sudo systemctl start semi-agentic.service

[Unit]
Description=Semi-Agentic Company Scheduler
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME_HERE
WorkingDirectory=/home/YOUR_USERNAME_HERE/semi-agentic-company
ExecStart=/home/YOUR_USERNAME_HERE/semi-agentic-company/venv/bin/python main_scheduler.py
Restart=always
RestartSec=10

# Environment variables (optional)
# Environment="PYTHONUNBUFFERED=1"

# Logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
